# deploy/cloudflared.service
# /etc/systemd/system/cloudflared.service
# Cloudflare 隧道守护进程，负责 HTTPS 边缘终止与子路径转发
# 安装：sudo cp cloudflared.service /etc/systemd/system/ && sudo systemctl daemon-reload && sudo systemctl enable --now cloudflared
# 获取 token：在 Cloudflare Zero Trust 控制台创建隧道后复制 token
[Unit]
Description=Cloudflare Tunnel (cloudflared)
After=network-online.target
Wants=network-online.target

[Service]
User=cloudflared
Group=cloudflared
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
TimeoutStartSec=0
Type=notify
# 隧道 token 通过环境文件注入，避免命令行暴露（ps/进程列表可见）
# 需创建 /etc/cloudflared/cloudflared.env，内容形如：TUNNEL_TOKEN=xxx，权限 0600
EnvironmentFile=/etc/cloudflared/cloudflared.env
ExecStart=/usr/local/bin/cloudflared tunnel --no-autoupdate run
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
